# package/kvmd-hid/kvmd-hid.service
[Unit]
Description=kvmd-hid
[Service]
Type=oneshot
ExecStart=/usr/bin/blikvm/enable-gadget.sh
[Install]
WantedBy=multi-user.target