# package/kvmd-web/kvmd-web.service
[Unit]
Description=BLIKVM - web controller daemon
After=network.target busybox-syslogd.service
Requires=network.target busybox-syslogd.service

[Service]
Type=simple
User=root
ExecStartPre=/bin/sleep 60
WorkingDirectory=/mnt/exec/release
ExecStart=/mnt/exec/release/server_app
RemainAfterExit=true
Restart=always

[Install]
WantedBy=local-fs.target